# docker-scripts/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mysql]
command=/usr/libexec/mysqld --user=mysql --datadir=/var/lib/mysql
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/mysql/mysql.log
stderr_logfile=/var/log/mysql/mysql-error.log

[program:wildfly]
command=/opt/wildfly/bin/standalone.sh -b 0.0.0.0 -bmanagement 0.0.0.0
user=wildfly
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/wildfly/wildfly.log
stderr_logfile=/var/log/wildfly/wildfly-error.log
environment=JAVA_HOME="/usr/lib/jvm/java-1.8.0-openjdk"

# Made with Bob
